# rootfs/etc/services.d/ble-scanner/run
#!/usr/bin/with-contenv bashio

cd /opt/ble_scanner
exec s6-setuidgid root python3 main.py 