# ble_scanner/rootfs/etc/services.d/ble-scanner/finish
#!/usr/bin/with-contenv bashio

bashio::log.info "BLE Scanner service stopped" 