# ble_scanner/rootfs/etc/services.d/ble-scanner/run
#!/usr/bin/with-contenv bashio

cd /opt/ble_scanner
exec s6-setuidgid root gunicorn --bind 0.0.0.0:8099 --workers 1 --timeout 120 main:app 