# ble_scanner/rootfs/etc/services.d/ble-scanner/run
#!/usr/bin/with-contenv bashio

cd /opt/ble_scanner
exec s6-setuidgid root gunicorn -c gunicorn.conf.py main:app 